# deploy/mt5bot.service
[Unit]
Description=MT5 trading bot
After=network.target

[Service]
Type=simple
WorkingDirectory={{REPO_PATH}}
ExecStartPre={{REPO_PATH}}/deploy/prestart.sh
ExecStart=/usr/bin/env python3 {{REPO_PATH}}/remote_api.py
Restart=on-failure
RestartSec=5
WatchdogSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
