# debian/mt5bot.service
[Unit]
Description=MT5 Bot Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mt5bot
ExecStart=/usr/bin/python3 -m mt5.generate_signals
Restart=on-failure

[Install]
WantedBy=multi-user.target
