# debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload || true
    systemctl enable mt5bot.service || true
fi

#DEBHELPER#
